# roles/network_shares/templates/mnt-tank.mount.j2
[Unit]
Description=\\neo\tank mount

[Mount]
What=//192.168.0.200/tank
Where=/mnt/tank
Type=cifs
Options=guest,uid={{ user_id }},gid={{ user_gid }},iocharset=utf8
