# roles/network_shares/files/mnt-tank.automount
[Unit]
Description=\\neo\tank automount

[Automount]
Where=/mnt/tank
TimeoutIdleSec=300

[Install]
WantedBy=multi-user.target
